# cookbooks/barito_market/templates/default/sidekiq.service.erb
[Unit]
Description=systemd script for running sidekiq for baritoMarket

[Service]
Restart=on-failure
RestartSec=2
StartLimitInterval=50
StartLimitBurst=10
ExecStart=/bin/bash -lc 'RAILS_ENV=<%= @env %> bundle exec sidekiq -C <%= @sidekiq_conf %>'
User=<%= @user %>
WorkingDirectory=<%= @app_directory %>

[Install]
WantedBy=multi-user.target
